# extra/runit/ircd/run
#!/bin/sh

setuser ircd /home/ircd/run/bin/ircd -foreground
